# systemd/taskwarrior-notify.service
[Unit]
Description=Taskwarrior due task notifications

[Service]
Type=oneshot
ExecStart=%h/.config/waybar/scripts/taskwarrior-notify.sh
